# udev/trackpoint.rules
SUBSYSTEM=="input", ATTR{name}=="*TrackPoint*", RUN+="/usr/bin/trackpoint"
